# hash_verify.cron
# cron file for hashi 
# probably run this at 6am every morning 
0 6 * * * root /etc/local/bin/hashii 
