# conf/hash_verify.cron
# cron file for hashi 
# probably run this at 6am every morning 
30 22 * * * root /etc/local/bin/hash_verify 
